# bot/softcat-radar.timer
[Unit]
Description=Run SOFT CAT Radar bot daily at 09:30

[Timer]
OnCalendar=*-*-* 09:30 Europe/London
Persistent=true

[Install]
WantedBy=timers.target
