# bot/softcat-prompts.timer
[Unit]
Description=Run SOFT CAT Prompt Library bot weekly on Wednesdays at 09:00

[Timer]
OnCalendar=Wed 09:00 Europe/London
Persistent=true

[Install]
WantedBy=timers.target
